# 01_build_image/build_assets/bootc-host-test.service
# Run the host startup test automatically at boot.
[Unit]
Description=Run bootc host startup test
After=network-online.target sshd.service nvidia-cdi-refresh.service devpod.service
Wants=network-online.target devpod.service

# Define how the host test runs.
[Service]
Type=oneshot
ExecStart=/opt/project/bootc_host_test.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

# Start this service during normal boot.
[Install]
WantedBy=multi-user.target
