# 01_build_image/build_assets/nvidia-cdi-refresh.path
[Unit]
Description=Watch for NVIDIA devices to refresh CDI

[Path]
PathExists=/dev/nvidiactl
Unit=nvidia-cdi-refresh.service

[Install]
WantedBy=multi-user.target
