# 01_build_image/build_assets/multi_tenant/agent_quadlet/agent-openclaw-runtime.container.tmpl
# OpenClaw runtime container for agent ${TENANT}/${AGENT}.
# Read-only filesystem, no privilege escalation. Mounts the per-tenant
# provisioner socket so the in-container agentctl can reach the host
# daemon, and the pod-local messaging-bridge directory so the runtime's
# verb-table router can read inbound message envelopes.
[Unit]
Description=OpenClaw runtime ${TENANT}/${AGENT}
After=${TENANT}-${AGENT}-pod.service
BindsTo=${TENANT}-${AGENT}-pod.service

[Container]
Image=${RUNTIME_IMAGE}
Pod=${TENANT}-${AGENT}.pod
Environment=OPENCLAW_TENANT=${TENANT}
Environment=OPENCLAW_AGENT=${AGENT}
Environment=OPENCLAW_NETWORK=${NETWORK}
Environment=OPENCLAW_MSGBUS_DIR=/run/messaging-bridge
Environment=OPENCLAW_MSGBUS_INBOUND=/run/messaging-bridge/inbound.sock
Environment=OPENCLAW_AGENTCTL_SOCKET=/run/agentctl/agentctl.sock
Volume=/run/openclaw-provisioner/tenants/${TENANT}.sock:/run/agentctl/agentctl.sock:ro,Z
Volume=${TENANT_HOME}/agents/${AGENT}/msgbus:/run/messaging-bridge:Z
ReadOnly=true
NoNewPrivileges=true

[Service]
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=default.target
