# 01_build_image/build_assets/multi_tenant/agent_quadlet/agent-credential-proxy.container.tmpl
# Credential proxy sidecar for agent ${TENANT}/${AGENT}.
# Identical to the onboarding-pod credential-proxy: mounts the tenant's
# broker socket from the host, exposes /run/credential-proxy/agent.sock
# pod-locally, refuses admin verbs at the proxy boundary.
[Unit]
Description=Credential proxy ${TENANT}/${AGENT}
After=${TENANT}-${AGENT}-pod.service
BindsTo=${TENANT}-${AGENT}-pod.service

[Container]
Image=quay.io/m0ranmcharles/fedora_init:credential-proxy
Pod=${TENANT}-${AGENT}.pod
Environment=OPENCLAW_TENANT=${TENANT}
Environment=OPENCLAW_BROKER_SOCKET=/run/credential-proxy/broker.sock
Environment=OPENCLAW_AGENT_SOCKET=/run/credential-proxy/agent.sock
Volume=/run/openclaw-broker/tenants/${TENANT}.sock:/run/credential-proxy/broker.sock:ro,Z
ReadOnly=true
NoNewPrivileges=true

[Service]
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=default.target
