# 01_build_image/build_assets/multi_tenant/agent_quadlet/agent-dev-env.container.tmpl
# Dev environment container for agent ${TENANT}/${AGENT}.
# Mutable inside the container; sudo allowed inside the container only.
# Tenant volumes are mounted under /workspace/<volume-name>.
[Unit]
Description=Dev env ${TENANT}/${AGENT}
After=${TENANT}-${AGENT}-pod.service
BindsTo=${TENANT}-${AGENT}-pod.service

[Container]
Image=${ENV_IMAGE}
Pod=${TENANT}-${AGENT}.pod
PodmanArgs=--device=nvidia.com/gpu=all
PodmanArgs=--security-opt=label=disable
Environment=OPENCLAW_TENANT=${TENANT}
Environment=OPENCLAW_AGENT=${AGENT}
${AGENT_VOLUMES_BLOCK}

[Service]
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=default.target
