# 01_build_image/build_assets/multi_tenant/openclaw-broker.service
[Unit]
Description=OpenClaw credential broker
Documentation=file:///usr/share/doc/openclaw/credential_broker.md
After=local-fs.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/openclaw-broker
Restart=on-failure
RestartSec=5s
RuntimeDirectory=openclaw-broker
RuntimeDirectoryMode=0755
StateDirectory=openclaw-platform/broker
StateDirectoryMode=0750
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
