# 01_build_image/build_assets/multi_tenant/tenant-openclaw-runtime.container.tmpl
# OpenClaw runtime container for an OpenClaw tenant.
# More locked down than the dev environment: read-only image filesystem,
# no host mounts, no privilege escalation. The per-tenant provisioner
# socket is bind-mounted so the agentctl CLI inside the container can
# self-provision new agent pods (Phase 3, see docs/concepts/agent_provisioning.md).
[Unit]
Description=OpenClaw runtime for tenant ${TENANT}
After=${TENANT}-onboard-pod.service
BindsTo=${TENANT}-onboard-pod.service

[Container]
Image=quay.io/m0ranmcharles/fedora_init:openclaw-runtime
Pod=${TENANT}-onboard.pod
Environment=OPENCLAW_TENANT=${TENANT}
Volume=/run/openclaw-provisioner/tenants/${TENANT}.sock:/run/agentctl/agentctl.sock:ro,Z
# The onboarding pod is SSH-driven, not message-driven: override the image
# default (Phase-4 router) and use the idle stub instead. Agent pods get
# the router; the onboarding pod stays a passive shell host.
Exec=/usr/local/bin/openclaw-runtime-stub.sh
ReadOnly=true
NoNewPrivileges=true

[Service]
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=default.target
