# 01_build_image/build_assets/multi_tenant/tenant-credential-proxy.container.tmpl
# Credential proxy sidecar for an OpenClaw tenant pod.
# Platform-controlled: pod-local UNIX socket interface, no broad mounts.
# Phase 2: forwards in-pod agent requests to the host openclaw-broker via
# a tenant-specific socket bind-mounted from the host.
[Unit]
Description=OpenClaw credential proxy for tenant ${TENANT}
After=${TENANT}-onboard-pod.service
BindsTo=${TENANT}-onboard-pod.service

[Container]
Image=quay.io/m0ranmcharles/fedora_init:credential-proxy
Pod=${TENANT}-onboard.pod
Environment=OPENCLAW_TENANT=${TENANT}
Environment=OPENCLAW_BROKER_SOCKET=/run/credential-proxy/broker.sock
Environment=OPENCLAW_AGENT_SOCKET=/run/credential-proxy/agent.sock
# Bind-mount the tenant-specific broker socket file from the host. The host
# socket is created by openclaw-broker.service and chowned to the tenant's
# UID/GID so the rootless container can reach it. Mount with :ro,Z so the
# tenant container cannot replace the file.
Volume=/run/openclaw-broker/tenants/${TENANT}.sock:/run/credential-proxy/broker.sock:ro,Z
ReadOnly=true
NoNewPrivileges=true

[Service]
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=default.target
